# LVM-crypt/crypt.service
# Put me in "/etc/systemd/system/"

[Unit]
Description=Mount Ecrypted Volumes

[Service]
ExecStart=/usr/local/bin/encrypt.sh
